# third-party/s3ql/usr/lib/systemd/system/s3ql-mounts.service
[Unit]
Description=Mount all s3ql volumes from /etc/s3ql/mounts
After=network-online.target
Wants=network-online.target
ConditionPathExists=/etc/s3ql/mounts

[Service]
Type=oneshot
ExecStart=/usr/sbin/s3ql-mount-all
ExecStop=/usr/sbin/s3ql-umount-all
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
