# third-party/nebula/DEBIAN/postrm
#!/bin/bash
set -euo pipefail
if [[ "$1" = "purge" ]]; then
	systemctl disable --now nebula.service || true
	userdel -r nebula 2>/dev/null || true
fi
